# scripts/ops/numo-feeds-ssm.service
# Numo price publisher on EC2 — secrets pulled from SSM at start (keys never on disk).
# Install: cp scripts/ops/numo-feeds-ssm.service /etc/systemd/system/numo-feeds.service
#          systemctl daemon-reload && systemctl enable --now numo-feeds
[Unit]
Description=Numo FX price publisher (Base mainnet, SSM secrets)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/risk-core
ExecStart=/home/ec2-user/risk-core/scripts/ops/run-with-ssm.sh /usr/bin/python3 scripts/publish_fx_feeds.py
Environment=PATH=/home/ec2-user/.foundry/bin:/usr/local/bin:/usr/bin:/bin
Environment=AWS_REGION=us-east-1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
